# nginx/conf/supervisor/000-php-fpm.conf
[program:php-fpm]
command=/usr/sbin/php-fpm7.4 -R -F -c /etc/php/7.4/fpm/php-fpm.conf
numprocs=1
autostart=true
autorestart=false
priority=5
redirect_stderr=true
stderr_logfile=/var/log/php-fpm.err.log
stdout_logfile=/var/log/php-fpm.out.log