# nginx/conf/supervisor/003-mariadb.conf
[program:mariadb]
command=/usr/bin/mysqld_safe >/dev/null 2>&1 &
numprocs=1
autostart=true
autorestart=false
priority=300
