# conf/supervisor/003-mariadb.conf
[program:mariadb]
command=/usr/sbin/mysqld
numprocs=1
autostart=true
autorestart=false
priority=300
