# conf/supervisor/003-mariadb.conf
[program:mariadb]
command=/usr/sbin/mysqld
numprocs=1
autostart=false
autorestart=false
priority=300
